# patches/rootfs/usbgadget-serial.service
[Unit]
Description=GiantBoard USB Gadget
Before=network.target

[Service]
Type=simple
ExecStart=/usr/bin/usbgadget-serial

[Install]
WantedBy=multi-user.target
